# sensord.service
[Unit]
Description=Savonia Measurements sensor poller

[Service]
Type=notify
ExecStart=/usr/sbin/samid

[Install]
WantedBy=multi-user.target